# rc.local
#!/bin/bash
(while [ "true" ]
do
        ivpn=$(nmcli con show uuid *** | grep VPN.VPN-STATE | awk '{print $2}')
        if [[ $ivpn != "5" ]]; then
                (sleep 1s && nmcli connection up uuid ***)
        else
                sleep 30s
        fi
done) &
